# rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bashio
# DaVinci Resolve Postgres addon — backend longrun.
#
# Waits up to ~10s for postgres' unix socket to appear, then exec's the
# aiohttp server. The wait avoids a noisy first-boot 500 from /api/admin/state
# if the backend wins the race over postgres' startup.
set -euo pipefail

bashio::log.info "starting backend (aiohttp on :8080); waiting for postgres socket"

SOCKET_DIR=/run/postgresql
for _ in $(seq 1 20); do
    if [ -S "${SOCKET_DIR}/.s.PGSQL.5432" ]; then
        bashio::log.info "postgres socket up; starting backend"
        break
    fi
    sleep 0.5
done

exec python3 /usr/local/bin/backend/server.py
